# src/static/mnt-pool1-media.automount
[Unit]
Description=Automount mnt-Pool1-Media

[Automount]
Where=/mnt/Pool1/Media

[Install]
WantedBy=multi-user.target